# gentoo/.config/systemd/user/dotfiles-sync.service
[Unit]
Description=Sync dotfiles to git

[Service]
Type=oneshot
ExecStart=/home/aarav/bin/dotfiles-sync
